# templates/default/cvat.service.erb
[Unit]
Description = OpenCV cvat Django server
After = syslog.target network.target remote-fs.target

[Service]
User = <%= node['cvat']['user'] %>
Restart=always
RestartSec=5
PIDFile = <%= node['cvat']['pid_file'] %>
# Prefixing a command with '-' means that you should still continue, even if it fails.
ExecStartPre=-/bin/chown <%= node['cvat']['user'] %> <%= node['cvat']['pid_file'] %>
ExecStart = <%= "#{#{node['conda']['dir']}/envs/cvat/bin/python /home/#{node['cvat']['user']}/cvat/manage.py runserver 0.0.0.0:8000" %>
ExecStop = /home/<%= "#{#{node['cvat']['user']}/cvat-stop.sh
Type = simple

[Install]
WantedBy = multi-user.target
