# templates/default/cvat.service.erb
[Unit]
Description = OpenCV cvat Django server
After = syslog.target network.target remote-fs.target

[Service]
User = <%= node['cvat']['user'] %>
Restart=always
RestartSec=5
Environment=JAVA_HOME=<%= @java_home %>
Environment=HADOOP_HOME=<%= node['hops']['base_dir']%>
Environment=HADOOP_CONF_DIR=<%= node['hops']['conf_dir']%>
Environment=HADOOP_USER_NAME=<%= node['hops']['hdfs']['user']%>
Environment=PROJECT=demo_deep_learning_admin000
#PIDFile = <%= node['cvat']['pid_file'] %>
# Prefixing a command with '-' means that you should still continue, even if it fails.
#ExecStartPre=-/bin/chown <%= node['cvat']['user'] %> <%= node['cvat']['pid_file'] %>
ExecStart = <%= "#{node['conda']['dir']}/envs/cvat/bin/python /home/#{node['cvat']['user']}/cvat/manage.py runserver 0.0.0.0:8000" %>
ExecStop = /home/<%= node['cvat']['user'] %>/cvat-stop.sh
Type = simple

[Install]
WantedBy = multi-user.target
